# crontabs/root
# DO NOT EDIT THIS FILE - edit the master and reinstall.
# (/tmp/crontab.DWgJsG/crontab installed on Sat Dec  3 14:57:55 2022)
# (Cron version -- $Id: crontab.c,v 2.13 1994/01/17 03:20:37 vixie Exp $)
* * * * * cd /project && php artisan schedule:run >> /dev/null 2>&1
